# systemd/amazon-packages-mqtt.timer
[Unit]
Description=Amazon Packages MQTT timer

[Timer]
OnCalendar=*:0/1
Unit=amazon-packages-mqtt.service

[Install]
WantedBy=timers.target
